# sqli-labs/mysql.ini
[program:mysql]
command=/bin/bash -c "/usr/libexec/mariadb-prepare-db-dir && /usr/bin/mysqld_safe --basedir=/usr"
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile=/tmp/supervisor-mysql.log
